# debian/knockdaemon2.postinst
#!/bin/bash

# Rc
echo "postinst: pass 2 complete"
update-rc.d -f knockdaemon2 defaults
echo "postinst: pass 2 complete"


echo "POSTINSTALL : Compile all py to pyc"
find /opt/knock/knockdaemon2/ -name '*.pyc' -delete
/opt/knock/knockdaemon2/bin/python -m compileall /opt/knock/knockdaemon2


# Restart
echo "postinst: pass 3"
service knockdaemon2 restart
echo "postinst: pass 3 complete"
